# crontab
*/1  *  *  *  *  root source /.env;/app/main.py
